# klipper-usb.rules
# Автоматическое монтирование USB для Klipper
SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]*", ACTION=="add", ENV{ID_FS_USAGE}=="filesystem", TAG+="systemd", RUN+="/usr/local/bin/klipper-usb-mount %k"
SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]*", ACTION=="remove", RUN+="/usr/local/bin/klipper-usb-mount %k"
